# .docker/base-with-services/s6-assets/s6-rc.d/postgresql/run
#!/command/execlineb -P

with-contenv

postgres -D /home/aiida/.postgresql
